# deploy/mastodon-bot.service
[Unit]
Description=Mastodon Follow/Like Bot (venv)
After=network.target

[Service]
Type=simple
User=michael
WorkingDirectory=/home/michael/mastodon-bot
EnvironmentFile=/home/michael/mastodon-bot/.env
ExecStart=/home/michael/mastodon-bot/.venv/bin/python /home/michael/mastodon-bot/mastodon_bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
